# api/deployment/crontab
0 0,12 * * * /bin/bash /src/blockstack/api/search/scripts/search_reindex.sh
0,12,24,36,48 * * * * /bin/bash /src/blockstack/api/search/scripts/search_update.sh