# api/deployment/crontab
0,20,40 * * * * /bin/bash /src/blockstack/api/search/scripts/search_update.sh
